# etc/cron.daily/restart-log2ram.sh
#!/bin/sh

systemctl restart log2ram.service
