# media/levien/f55adb7f-6d96-4e27-a51d-a274dfa58214/etc/udev/rules.d/10-Huawei.rules
SUBSYSTEMS=="usb", ATTRS{modalias}=="usb:v12D1p1F01*", SYMLINK+="hwcdrom", RUN+="/usr/bin/sg_raw /dev/hwcdrom 11 06 20 00 00 00 00 00 01 00"
